# infra/docker/local/consumer/etc/supervisor/conf.d/messenger-worker.conf
[program:messenger-consume]
command=/home/app/bin/console messenger:consume --time-limit=3600
user=root
numprocs=1
startsecs=0
autostart=true
autorestart=true
process_name=%(program_name)s_%(process_num)02d
stderr_logfile=/var/log/supervisor/%(program_name)s_stderr.log
stderr_logfile_maxbytes=10MB
stdout_logfile=/var/log/supervisor/%(program_name)s_stdout.log
stdout_logfile_maxbytes=10MB
